# config/dev/cloud/ap/ui/crudusrs/config/proxy/renewer/crontab
# min	hour	day	month	weekday	command
0 0 * * MON step ca renew --force /var/local/step/site.crt /var/local/step/site.key
